# jehon-system-private/etc/cron.daily/jh-synology-torrents
#!/usr/bin/env bash

if [ -r /home/jehon/Documents/src/packages/bin/jh-synology-torrents.sh ]; then
    runuser -u jehon -- /home/jehon/Documents/src/packages/bin/jh-synology-torrents.sh
fi
